# scripts/com.aimarket.railway-config-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.aimarket.railway-config-export</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/max/Projects/ai-market/runbooks/scripts/run_railway_config_export.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>EnvironmentVariables</key>
    <dict><key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
    <key>StandardOutPath</key><string>/Users/max/Library/Logs/aimarket_railway_config_export.out.log</string>
    <key>StandardErrorPath</key><string>/Users/max/Library/Logs/aimarket_railway_config_export.error.log</string>
    <key>RunAtLoad</key><false/>
    <key>Nice</key><integer>10</integer>
</dict>
</plist>
